# systemd/.config/systemd/user/swaync.service
[Unit]
Description=SwayNC Service
After=graphical.target

[Service]
Type=simple
ExecStart=swaync
restart=always

[Install]
WantedBy=default.target
