# systemd/.config/systemd/user/waybar.service
[Unit]
Description=Waybar Service
After=graphical.target

[Service]
Type=simple
ExecStart=waybar
restart=always

[Install]
WantedBy=default.target
